# macromilter/MacroMilter.conf
description "MacroMilter Service"

start on runlevel [234]
stop on runlevel [0156]

setuid postfix
setgid postfix

chdir /etc/macromilter
exec /usr/bin/python macromilter.py > run.log
respawn
