# macromilter/MacroMilter.conf
description "MacroMilter Service"

start on runlevel [234]
stop on runlevel [0156]

setuid postfix
setgid postfix

chdir /usr/bin/
exec /usr/bin/python macromilter.py
respawn
